# supervisor/hc/conf.d/music_idea.conf
[program:music_idea]
command = /home/utylee/.virtualenvs/misc/bin/python /home/utylee/.virtualenvs/misc/src/music-idea.py

autostart = true
autorestart = true

