# supervisor/hc/conf.d/weather.conf
[program:weather]
command = /home/utylee/.virtualenvs/misc/bin/python /home/utylee/.virtualenvs/misc/src/weather.py

autostart = true
autorestart = true

